# systemd/user/theme-daemon.service
[Unit]
Description=Apply N4P1x dynamic theme at session start
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/theme/sync-all
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
